# backend-laravel/docker/supervisord.conf
[supervisord]
nodaemon=true

[unix_http_server]
file=/tmp/supervisor.sock
chmod=0700

[program:php-fpm]
command=/usr/local/sbin/php-fpm -F
autorestart=true
priority=10

[program:nginx]
command=/usr/sbin/nginx -g 'daemon off;'
autorestart=true
priority=20
